# deploy/yoji.service
[Unit]
Description=Yoji portfolio site
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/yoji
ExecStart=/usr/bin/env node /var/www/yoji/server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=7878

[Install]
WantedBy=multi-user.target
